# scraper/scraper.service
[Unit]
Description=Execute scraper.sh

[Service]
ExecStart=/opt/scraper/scraper.sh
Type=oneshot
User=scraper

[Install]
WantedBy=multi-user.target