# mongodb/backup.timer
[Unit]
Description=launch the backup script

[Timer]
OnCalendar=*-*-* 5:00:00
Persistent=true

[Install]
WantedBy=timers.target

